# root/etc/hotplug.d/iface/98-bypass
#!/bin/sh
# hotplug.d/iface/98-bypass
#
# Keep firewall WAN sets fresh. Changes to any explicitly selected Naive or
# Direct egress interface require a restart so netifd's new L3 device/gateway
# or BypassCore device bind is applied.

[ -f /var/lock/bypass_ready.lock ] || exit 0

. /lib/functions/network.sh
. /usr/share/bypass/utils.sh

[ "$(uci -q get bypass.@global[0].enabled)" = "1" ] || exit 0

# Guard against re-entry.
LOCK=/var/lock/bypass_ifup.lock
exec 9>"$LOCK"
flock -xn 9 || exit 0

direct_iface=$(uci -q get bypass.@global_rules[0].direct_egress_interface)
node_egress_ifaces=$(
	for sid in $(uci -q show bypass 2>/dev/null | sed -n 's/^bypass\.\([^.=]*\)=shunt_rules$/\1/p'); do
		node=$(uci -q get bypass.${sid}.outbound)
		[ "$(uci -q get bypass.${node})" = "nodes" ] || continue
		uci -q get bypass.${node}.egress_interface
	done | awk 'NF && !seen[$0]++'
)
rule_direct_ifaces=$(
	for sid in $(uci -q show bypass 2>/dev/null | sed -n 's/^bypass\.\([^.=]*\)=shunt_rules$/\1/p'); do
		[ "$(uci -q get bypass.${sid}.outbound)" = "_direct" ] || continue
		uci -q get bypass.${sid}.egress_interface
	done | awk 'NF && !seen[$0]++'
)
is_managed_iface() {
	[ -n "$INTERFACE" ] || return 1
	[ "$INTERFACE" = "$direct_iface" ] && return 0
	case " $node_egress_ifaces " in *" $INTERFACE "*) return 0 ;; esac
	case " $rule_direct_ifaces " in *" $INTERFACE "*) return 0 ;; esac
	return 1
}

is_default_route_iface() {
	# DEVICE from hotplug; check it carries the default route.
	local dev="$DEVICE"
	[ -z "$dev" ] && return 1
	ip route show default 2>/dev/null | grep -q "dev $dev"
}

case "$ACTION" in
	ifupdate)
		if is_managed_iface; then
			( /etc/init.d/bypass restart ) >/dev/null 2>&1 &
		else
			use_tables=$(get_cache_var USE_TABLES)
			[ -n "$use_tables" ] && [ -x "$APP_PATH/${use_tables}.sh" ] && \
				"$APP_PATH/${use_tables}.sh" update_wan_sets 2>/dev/null
		fi
		;;
	ifup)
		if is_default_route_iface || is_managed_iface; then
			( /etc/init.d/bypass restart ) >/dev/null 2>&1 &
		fi
		;;
	ifdown)
		# Any selected Naive or Direct WAN disappearing must remove its stale
		# route/bind immediately instead of silently falling back to another WAN.
		is_managed_iface && \
			( /etc/init.d/bypass restart ) >/dev/null 2>&1 &
		;;
esac

flock -u 9
exit 0
