# root/etc/hotplug.d/iface/98-bypass
#!/bin/sh
# hotplug.d/iface/98-bypass
#
# Keep firewall WAN sets fresh. Changes to any explicitly selected Naive or
# Direct egress interface require a restart so netifd's new L3 device/gateway
# or BypassCore device bind is applied.

. /lib/functions/network.sh
. /usr/share/bypass/utils.sh

[ "$(uci -q get bypass.@global[0].enabled)" = "1" ] || exit 0

# If cold-boot startup exhausted its bounded retries before the WAN became
# usable, a later ifup must be able to recover the service. Do not interfere
# with the configured boot delay while its marker exists.
if [ ! -f /var/lock/bypass_ready.lock ]; then
	if { [ "$ACTION" = "ifup" ] || [ "$ACTION" = "ifupdate" ]; } && \
	   [ ! -f /var/lock/bypass_booting.lock ]; then
		if mkdir /var/lock/bypass_recovering 2>/dev/null; then
			( /etc/init.d/bypass restart; rmdir /var/lock/bypass_recovering 2>/dev/null ) >/dev/null 2>&1 &
		fi
	fi
	exit 0
fi

# Guard against re-entry.
LOCK=/var/lock/bypass_ifup.lock
exec 9>"$LOCK"
flock -xn 9 || exit 0

direct_iface=$(uci -q get bypass.@global_rules[0].direct_egress_interface)
default_naive_iface=$(uci -q get bypass.@global_rules[0].default_naive_interface)
node_egress_ifaces=$(
	{
		default_node=$(uci -q get bypass.@global_rules[0].default_node)
		if [ "$(uci -q get bypass.${default_node})" = "nodes" ]; then
			[ "$(uci -q get bypass.${default_node}.node_type)" = "wireguard" ] && default_node=""
		fi
		if [ -n "$default_node" ] && [ "$(uci -q get bypass.${default_node})" = "nodes" ]; then
			iface=$(uci -q get bypass.${default_node}.egress_interface)
			printf '%s\n' "${iface:-$default_naive_iface}"
		fi
		for sid in $(uci -q show bypass 2>/dev/null | sed -n 's/^bypass\.\([^.=]*\)=shunt_rules$/\1/p'); do
			[ "$(uci -q get bypass.${sid}.is_default)" = "1" ] && continue
			node=$(uci -q get bypass.${sid}.outbound)
			[ "$(uci -q get bypass.${node})" = "nodes" ] || continue
			[ "$(uci -q get bypass.${node}.node_type)" = "wireguard" ] && continue
			iface=$(uci -q get bypass.${node}.egress_interface)
			printf '%s\n' "${iface:-$default_naive_iface}"
		done
	} | awk 'NF && !seen[$0]++'
)
rule_direct_ifaces=$(
	for sid in $(uci -q show bypass 2>/dev/null | sed -n 's/^bypass\.\([^.=]*\)=shunt_rules$/\1/p'); do
		[ "$(uci -q get bypass.${sid}.is_default)" = "1" ] && continue
		[ "$(uci -q get bypass.${sid}.outbound)" = "_direct" ] || continue
		uci -q get bypass.${sid}.egress_interface
	done | awk 'NF && !seen[$0]++'
)
is_managed_iface() {
	[ -n "$INTERFACE" ] || return 1
	[ "$INTERFACE" = "$direct_iface" ] && return 0
	case " $node_egress_ifaces " in *" $INTERFACE "*) return 0 ;; esac
	case " $rule_direct_ifaces " in *" $INTERFACE "*) return 0 ;; esac
	return 1
}

is_direct_iface() {
	[ -n "$INTERFACE" ] || return 1
	[ "$INTERFACE" = "$direct_iface" ] && return 0
	case " $rule_direct_ifaces " in *" $INTERFACE "*) return 0 ;; esac
	return 1
}

is_default_route_iface() {
	# DEVICE from hotplug; check it carries the default route.
	local dev="$DEVICE"
	[ -z "$dev" ] && return 1
	ip route show default 2>/dev/null | grep -q "dev $dev"
}

managed_runtime_matches() {
	# A routine DHCP renew commonly emits ifupdate without changing device or
	# gateway. Passwall2 only refreshes WAN sets for that event. For our
	# per-node policy tables, also verify the active table before deciding that a
	# full restart is needed.
	get_egress_runtime "$INTERFACE" || return 1
	local table iface route route6
	if [ -s "$TMP_PATH/egress_tables" ]; then
		while read -r table iface; do
			[ "$iface" = "$INTERFACE" ] || continue
			route=$(ip -o route show table "$table" default proto 99 2>/dev/null)
			route6=$(ip -6 -o route show table "$table" default proto 99 2>/dev/null)
			[ -n "$route$route6" ] || return 1
			if [ -n "$route" ]; then
				if [ -n "$EGRESS_GATEWAY" ]; then
					echo "$route" | grep -q "via $EGRESS_GATEWAY .*dev $EGRESS_DEVICE\( \|$\)" || return 1
				else
					echo "$route" | grep -q "dev $EGRESS_DEVICE\( \|$\)" || return 1
				fi
			fi
			if [ -n "$route6" ]; then
				if [ -n "$EGRESS_GATEWAY6" ]; then
					echo "$route6" | grep -q "via $EGRESS_GATEWAY6 .*dev $EGRESS_DEVICE\( \|$\)" || return 1
				else
					echo "$route6" | grep -q "dev $EGRESS_DEVICE\( \|$\)" || return 1
				fi
			fi
		done < "$TMP_PATH/egress_tables"
	fi

	# Direct outbounds bind to the L3 device, not its renewed address. If this
	# logical interface is used for Direct, confirm the generated core config
	# still names the current device.
	if is_direct_iface; then
		jsonfilter -i "$BYPASSCORE_CFG" -e '@.outbounds[*].bind.interface' 2>/dev/null | \
			grep -Fxq "$EGRESS_DEVICE" || return 1
	fi
	return 0
}

case "$ACTION" in
	ifupdate)
		if is_managed_iface; then
			if ! managed_runtime_matches; then
				( /etc/init.d/bypass restart ) >/dev/null 2>&1 &
			fi
		fi
		use_tables=$(get_cache_var USE_TABLES)
		[ -n "$use_tables" ] && [ -x "$APP_PATH/${use_tables}.sh" ] && \
			"$APP_PATH/${use_tables}.sh" update_wan_sets 2>/dev/null
		;;
	ifup)
		if is_default_route_iface || is_managed_iface; then
			( /etc/init.d/bypass restart ) >/dev/null 2>&1 &
		fi
		;;
	ifdown)
		# Any selected Naive or Direct WAN disappearing must remove its stale
		# route/bind immediately instead of silently falling back to another WAN.
		is_managed_iface && \
			( /etc/init.d/bypass restart ) >/dev/null 2>&1 &
		;;
esac

flock -u 9
exit 0
