# root/etc/init.d/bypass
#!/bin/sh /etc/rc.common

START=99
STOP=15

. $IPKG_INSTROOT/usr/share/bypass/utils.sh
APP_FILE=${APP_PATH}/app.sh
LOCK_FILE_DIR=/var/lock
LOCK_FILE=${LOCK_FILE_DIR}/${CONFIG}.lock

set_lock() {
	[ ! -d "$LOCK_FILE_DIR" ] && mkdir -p "$LOCK_FILE_DIR"
	exec 999>"$LOCK_FILE"
	flock -x -w 10 999
}

unset_lock() {
	flock -u 999 2>/dev/null
}

boot_func() {
	local delay
	delay=$(uci -q get ${CONFIG}.@global_delay[0].start_delay 2>/dev/null || echo 1)
	[ "$delay" -gt 0 ] 2>/dev/null && {
		log 0 "Start after a delay of %s seconds." "$delay"
		sleep "$delay"
	}
	restart boot
	touch ${LOCK_FILE_DIR}/${CONFIG}_ready.lock
}

boot() {
	boot_func >/dev/null 2>&1 &
}

start() {
	set_lock || { log 0 "Could not acquire service lock, exit."; exit 1; }
	$APP_FILE start
	local rc=$?
	unset_lock
	return $rc
}

stop() {
	set_lock || { log 0 "Could not acquire service lock for stop, exit."; exit 1; }
	$APP_FILE stop
	unset_lock
}

restart() {
	set_lock || { log 0 "Could not acquire service lock for restart, exit."; exit 1; }
	$APP_FILE stop
	$APP_FILE start $1
	local rc=$?
	unset_lock
	return $rc
}

reload() {
	echo "Warning: this service does not support live reload; full restart instead."
	restart
}
